# packaging/com.nyx.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nyx.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/nyx-agent</string>
        <string>serve</string>
        <string>--headless</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>/usr/local/var/lib/nyx-agent</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NYX_AGENT_STATE_DIR</key>
        <string>/usr/local/var/lib/nyx-agent</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/usr/local/var/log/nyx-agent.out.log</string>

    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/nyx-agent.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
