# packaging/nyctos.timer
[Unit]
Description=Trigger a weekly nyctos scan
Documentation=https://nyctos.dev/docs/triggers/

[Timer]
# 03:00 every Monday — matches the Phase 27 acceptance fixture
# (`0 3 * * 1`). Operators can override by dropping a unit file at
# /etc/systemd/system/nyctos.timer.d/override.conf.
OnCalendar=Mon *-*-* 03:00:00
Persistent=true
Unit=nyctos-scan.service

[Install]
WantedBy=timers.target
