# packaging/nyx-agent.service
[Unit]
Description=Nyx Agent scanning daemon
Documentation=https://nyxscan.dev/docs/triggers/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# The systemd timer in `nyx-agent.timer` triggers a one-shot scan via
# the daemon's `POST /api/v1/scan` endpoint, but the daemon itself
# should run continuously so the API + WebSocket surfaces stay
# reachable and the in-process cron scheduler can fire `[[schedule]]`
# entries. Operators who prefer to drive scans purely via the timer
# can set Type=oneshot and replace ExecStart with the scan command.
ExecStart=/usr/local/bin/nyx-agent --state-dir %S/nyx-agent serve --headless
Restart=on-failure
RestartSec=5

# Hardening — the daemon does not need root, network admin, or write
# access to anything outside its state directory.
DynamicUser=yes
StateDirectory=nyx-agent
WorkingDirectory=%S/nyx-agent
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
SystemCallArchitectures=native
# AmbientCapabilities/CapabilityBoundingSet intentionally empty.

[Install]
WantedBy=multi-user.target
